# deploy/systemd/galliasphere.service
[Unit]
Description=update galliasphere repository

[Service]
WorkingDirectory=/home/core/repository/galliasphere
ExecStart=/usr/bin/git pull
